# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
silent=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:tor]
command=/usr/bin/tor --DataDirectory /var/lib/tor --SocksPort 0.0.0.0:9050 --ControlPort 0.0.0.0:9051 --SocksPolicy "accept 10.0.0.0/8" --SocksPolicy "accept 172.16.0.0/12" --SocksPolicy "accept 192.168.0.0/16" --SocksPolicy "reject *" --Log "notice stdout" --CircuitBuildTimeout 5 --KeepalivePeriod 60 --NewCircuitPeriod 15 --NumEntryGuards 8
autostart=true
autorestart=true
priority=200
stdout_logfile=/tmp/tor.log
stderr_logfile=/tmp/tor.log
user=tor